# systemd/.config/systemd/user/input-leap.service
[Unit]
Description=Input Leap Client
After=network.target

[Service]
ExecStart=input-leap
Restart=always

[Install]
WantedBy=default.target
